# infra/systemd/gept-validation.timer
[Unit]
Description=Run prediction validation hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
