# infra/systemd/gept-ml-validation.timer
[Unit]
Description=Run ML ranker validation daily

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
