# infra/systemd/gept-collect-hourly.service
[Unit]
Description=GePT hourly data collector
After=network.target postgresql.service

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/gept/current/packages/model
Environment="PATH=/home/ubuntu/gept/venv/bin:/usr/bin"
EnvironmentFile=/home/ubuntu/gept/.env
ExecStart=/home/ubuntu/gept/venv/bin/python collectors/collect_hourly_pg.py
TimeoutStartSec=300
