# infra/systemd/gept-trainer.service
[Unit]
Description=GePT Model Trainer - Opportunistic Validation-Driven Training
After=network.target tailscaled.service
Wants=tailscaled.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/gept
Environment="PATH=/home/ubuntu/gept/venv/bin:/usr/bin"
Environment="PYTHONPATH=/home/ubuntu/gept/src:/home/ubuntu/gept"
EnvironmentFile=/home/ubuntu/gept/.env
ExecStart=/home/ubuntu/gept/venv/bin/python -m src.gept_trainer
Restart=always
RestartSec=60

# Graceful shutdown - give it time to finish current model
TimeoutStopSec=600
KillSignal=SIGTERM

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gept-trainer

[Install]
WantedBy=multi-user.target
