# infra/systemd/gept-engine.service
[Unit]
Description=GePT Recommendation Engine API (Docker)
After=network.target docker.service
Wants=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=ubuntu
WorkingDirectory=/home/ubuntu/gept/current/packages/engine
Environment="PATH=/usr/bin:/bin"
# Always pull before starting so a `systemctl restart gept-engine` picks up the
# latest image (or the tag specified by GEPT_ENGINE_TAG).
ExecStartPre=/usr/bin/docker compose -f /home/ubuntu/gept/current/packages/engine/docker-compose.yml pull
ExecStart=/usr/bin/docker compose -f /home/ubuntu/gept/current/packages/engine/docker-compose.yml up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f /home/ubuntu/gept/current/packages/engine/docker-compose.yml down
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
