# packages/model/config/systemd/gept-validation.timer
[Unit]
Description=Run GePT Prediction Validation Hourly

[Timer]
# Run every hour at minute 30 (offset from inference at :00/:05)
OnCalendar=*:30
Persistent=true

# Randomize start slightly to avoid thundering herd
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
