# packages/model/config/systemd/gept-validation.service
[Unit]
Description=GePT Prediction Validation Job
After=network.target postgresql.service

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu

# Working directory
WorkingDirectory=/home/ubuntu/gept/current/packages/model

# Environment
Environment=DATABASE_URL=postgresql://gept_app@localhost/gept

# Run the validation job
ExecStart=/home/ubuntu/gept/venv/bin/python run_validation.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gept-validation

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
